# etc/profile.d/xdg.sh
# set xdg config path
export XDG_CONFIG_HOME="/home/alexandre/.config"
